# ansible/roles/dns/templates/pihole.service.j2
[Unit]
Description=Pod container-pihole.service
Wants=network.target
After=network-online.target

[Service]
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=kubectl delete -f "/home/{{ user_name }}/kubernetes/pihole.yml"
ExecStart=kubectl create -f "/home/{{ user_name }}/kubernetes/pihole.yml"
ExecStop=kubectl delete -f "/home/{{ user_name }}/kubernetes/pihole.yml"
Type=forking

[Install]
WantedBy=multi-user.target default.target
