# roles/nextcloud/templates/nextcloud_container.service.j2
[Unit]
Description=Pod container-nextcloud.service
Wants=network.target
After=network-online.target

[Service]
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=kubectl delete -f "/home/{{ user_name }}/kubernetes/nextcloud.yml"
ExecStart=kubectl create -f "/home/{{ user_name }}/kubernetes/nextcloud.yml"
ExecStop=kubectl delete -f "/home/{{ user_name }}/kubernetes/nextcloud.yml"
Type=forking

[Install]
WantedBy=multi-user.target default.target
