# roles/navidrome/templates/check_music_rclone.timer.j2
[Unit]
Description=Timer for check_music_rclone.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=check_music_rclone.service

[Install]
WantedBy=timers.target
